# launchagents/com.claude.obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.obsidian-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/takamichi.omori/workspace/github.com/onsd/dotfiles/script/watch-and-save.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/takamichi.omori/workspace/github.com/onsd/dotfiles/logs/obsidian-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/takamichi.omori/workspace/github.com/onsd/dotfiles/logs/obsidian-sync-error.log</string>
</dict>
</plist>
